# autostart/picom.desktop
[Desktop Entry]
Type=Application
Name=polybar
Exec=picom
StartupNotify=false
Terminal=false
