# autostart/unclutter.desktop
[Desktop Entry]
Type=Application
Name=unclutter
Exec=unclutter -idle 4 -root
StartupNotify=false
Terminal=false
